# checkusb.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Polling of checking usb"
author          "Sakya   Sakyasir@163.com"

start on started boot-services
  
script

if [ -e /usr/share/cros/check_usb.sh ]; then 

. /usr/share/cros/check_usb.sh

fi

end script

