# mybluetooth.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the bluetooth"
author          "Sakyasir@163.com"

start on started boot-services
stop on pre-shutdown

respawn

script

  exec /usr/sbin/rtk_hciattach -n -s 115200 ttyS1 rtk_h5

end script
